# .win/wsl/distro/ubuntu/services/user/obsidian-sync.service
[Unit]
Description=Obsidian Vault Auto Sync Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'export HOME=$(getent passwd $(id -un) | cut -d: -f6); /usr/local/bin/obsidian-sync'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 環境変数を明示的に設定
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/sbin

# ワーキングディレクトリはランタイムで決定
# WorkingDirectory=%h

# プロセスの制限
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
